# misc/systemd/asus-fan.service
[Unit]                                                                                        
Description=Create tmp-files to provide a consistent naming and location

[Service]                                                                                     
ExecStart=/usr/bin/asus-fan-create-symlinks.sh                      

[Install]                                                                                     
WantedBy=multi-user.target
Before=thermald.service
